# mkosi.skeleton/etc/systemd/system/confer-docling.service
[Unit]
Description=Confer Docling Service
After=local-fs.target mnt-config.mount nvidia-persistenced.service nvidia-cc-attestation.service
Requires=mnt-config.mount nvidia-cc-attestation.service

[Service]
Type=simple
User=root
EnvironmentFile=/mnt/config/secrets.env

# Ensure artifacts directory exists
ExecStartPre=/bin/mkdir -p /tmp/docling-models

# Run with clean environment - only pass specific variables from secrets.env
ExecStart=/bin/bash -c 'exec env -i \
    PATH=/opt/venv-docling/bin:/usr/local/bin:/usr/bin:/bin \
    HOME=/tmp \
    XDG_CACHE_HOME=/tmp/huggingface \
    DOCLING_SERVE_ARTIFACTS_PATH=/tmp/docling-models \
    CUDA_VISIBLE_DEVICES=0 \
    DOCLING_DEVICE="${DOCLING_DEVICE:-cuda}" \
    DOCLING_DEBUG_PROFILE_PIPELINE_TIMINGS=true \
    DOCLING_SERVE_ENG_LOC_NUM_WORKERS="${DOCLING_SERVE_ENG_LOC_NUM_WORKERS:-2}" \
    DOCLING_SERVE_ENG_LOC_SHARE_MODELS="${DOCLING_SERVE_ENG_LOC_SHARE_MODELS:-true}" \
    DOCLING_NUM_THREADS="${DOCLING_NUM_THREADS:-4}" \
    DOCLING_PERF_PAGE_BATCH_SIZE="${DOCLING_PERF_PAGE_BATCH_SIZE:-4}" \
    DOCLING_PERF_ELEMENTS_BATCH_SIZE="${DOCLING_PERF_ELEMENTS_BATCH_SIZE:-8}" \
    DOCLING_SERVE_OCR_BATCH_SIZE="${DOCLING_SERVE_OCR_BATCH_SIZE:-}" \
    DOCLING_SERVE_LAYOUT_BATCH_SIZE="${DOCLING_SERVE_LAYOUT_BATCH_SIZE:-}" \
    DOCLING_SERVE_TABLE_BATCH_SIZE="${DOCLING_SERVE_TABLE_BATCH_SIZE:-}" \
    /opt/venv-docling/bin/docling-serve run --host 127.0.0.1 --port 5001'

Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
